# imu_ultrasonic_data/10-ftdi-he.rules
    ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6001", ATTRS{serial}=="A5022WCW", RUN+="/bin/sh -c 'echo $kernel > /sys/bus/usb/drivers/ftdi_sio/unbind'"
    ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6001", ATTRS{serial}=="A5022WCW", RUN+="/bin/sh -c 'echo $kernel > /sys/bus/usb/drivers/usbserial/unbind'"
    SUBSYSTEM=="usb", ATTRS{idVendor}=="0403",GROUP="dialout"